# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
logfile_backups=0
loglevel=info
pidfile=/var/run/supervisord.pid
user=root

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:uvicorn]
command=bash -c '/opt/venv/bin/python -m uvicorn backend.main:app --host 0.0.0.0 --port ${WEB_PORT:-8080} --workers 1 --no-access-log'
directory=/app
autostart=true
autorestart=true
startretries=5
startsecs=5
stopsignal=TERM
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stdout_logfile_backups=0
environment=PYTHONPATH="/app",PYTHONUNBUFFERED="1"
